# root/etc/services.d/beanstalk/run
#!/bin/bash

s6-setuidgid alpine exec beanstalkd

