# root/etc/services.d/kanboard-worker/run
#!/bin/bash

cd /config/www/kanboard || exit;

s6-setuidgid alpine exec ./cli worker


